# lib/launchd/com.clawkey.proxy.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Clawkey LiteLLM proxy — user-scope launchd agent.
  Installed to ~/Library/LaunchAgents/com.clawkey.proxy.plist by `clawkey proxy install`.
  Domain: gui/$UID. Requires no admin privileges.

  Placeholders (replaced by `clawkey proxy install`):
    {{CLAWKEY_DIR}}        absolute path to the clawkey checkout
    {{CLAWKEY_CONFIG_DIR}} XDG config dir (~/.config/clawkey by default)
    {{CLAWKEY_STATE_DIR}}  XDG state dir (~/.local/state/clawkey by default)
    {{LITELLM_PORT}}       port number (default 4040)
    {{LOG_PATH}}           absolute path to the proxy log file
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clawkey.proxy</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{CLAWKEY_DIR}}/lib/launchd/run-proxy.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>CLAWKEY_DIR</key>
        <string>{{CLAWKEY_DIR}}</string>
        <key>CLAWKEY_CONFIG_DIR</key>
        <string>{{CLAWKEY_CONFIG_DIR}}</string>
        <key>CLAWKEY_STATE_DIR</key>
        <string>{{CLAWKEY_STATE_DIR}}</string>
        <key>LITELLM_PORT</key>
        <string>{{LITELLM_PORT}}</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>{{CLAWKEY_DIR}}</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>{{LOG_PATH}}</string>

    <key>StandardErrorPath</key>
    <string>{{LOG_PATH}}</string>
</dict>
</plist>
